# supervisord.conf
[program:resque-worker]
command = bundle exec rake resque:work
environment = QUEUE=*
directory = /home/dotcloud/current
